# assets/supervisord.aptlyapi.conf
[program:aptlyapi]
command=/usr/bin/aptly api serve -no-lock
redirect_stderr=true
stdout_logfile=/opt/logs/aptlyapi.log